# scripts/gunicorn.socket
# /etc/systemd/system/gunicorn.socket
#
# Socket-activation keeps the Unix socket open across gunicorn reloads.
# nginx connects to this socket; requests queue here while workers are
# being gracefully replaced, so no connections are refused during deploys.

[Unit]
Description=Gunicorn socket for Fitness AI App

[Socket]
ListenStream=/run/gunicorn.sock
# The socket is owned by www-data so both nginx and gunicorn can use it
SocketUser=www-data

[Install]
WantedBy=sockets.target
