# scripts/gunicorn.service
# /etc/systemd/system/gunicorn.service
#
# Notes on the flags chosen here and why they matter for deployments:
#
#   --workers 3         Standard formula: 2 * CPU cores + 1.  Adjust to your Droplet size.
#   --preload           Loads app code in the master *before* forking workers.
#                       On reload (SIGHUP), new workers inherit the already-loaded
#                       code, so they start up faster and the reload window is shorter.
#   --timeout 120       Worker is killed if it doesn't respond within 120 s.
#                       Prevents stuck workers from blocking slots indefinitely.
#   --graceful-timeout 30
#                       After a SIGHUP the master waits up to 30 s for old workers
#                       to finish their current request before sending SIGKILL.
#                       Set this higher than your longest expected request time.
#   --keep-alive 5      Reuse connections for 5 s – reduces overhead from nginx upstream.

[Unit]
Description=Gunicorn daemon for Fitness AI App
# Start after the socket unit so socket-activation works correctly
Requires=gunicorn.socket
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/fitness-ai-app/fitness_ai_app

ExecStart=/var/www/fitness-ai-app/venv/bin/gunicorn \
    --workers 3 \
    --bind unix:/run/gunicorn.sock \
    --timeout 120 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --preload \
    --access-logfile /var/log/gunicorn/access.log \
    --error-logfile  /var/log/gunicorn/error.log \
    fitness_ai_app.wsgi:application

# SIGHUP triggers a graceful in-place reload (new workers, old finish requests)
ExecReload=/bin/kill -HUP $MAINPID

# Restart the service automatically if it crashes, but not on clean stops
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
